# daemon/claude-meter.service
[Unit]
Description=Claude Usage Tracker Daemon (BLE)
Requires=bluetooth.target
After=bluetooth.target

[Service]
Type=simple
WorkingDirectory=REPO_DIR
ExecStart=EXEC_CMD
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
